# clcbio-server-sdk/conf/clcserver-startupscript
#! /bin/sh 
# 
# CLC Genomics Server startup script
# Copyright (C) 2002-2011 CLC bio
# 
# chkconfig: - 99 04
# description: CLC Genomics Server
  
# Set the username under which the  server will run
USER=clcserver 

# Set the path to the directory that contains the server files 
DIRECTORY=/usr/local/CLCGenomicsServer 

# Set the path to executable
EXEC=$DIRECTORY/CLCGenomicsServer
# 
#---------------------------------------------------------------- 
# DO NOT EDIT BELOW THIS LINE 
#---------------------------------------------------------------- 
# 

test -x $EXEC || exit 0

start() { 
echo $debuglog 
 su - $USER -c "$EXEC start" 
} 

stop() { 
 su - $USER -c "$EXEC stop" 
} 

status() { 
 su - $USER -c "$EXEC status" 
} 

case "$1" in 
 start) 
	start 
 ;; 
 stop) 
	stop 
 ;;
 status) 
	status
 ;; 
 restart) 
	stop 
	sleep 2 
	start 
 ;; 
 *) 
 echo $"Usage: $0 {start|stop|restart|status}" 
 exit 1 
esac 

exit 0
